# pi-hole/rootfs/etc/cron.d/unbound
# Unbound: A validating, recursive, and caching DNS resolver. 
# https://nlnetlabs.nl/projects/unbound/about/
# https://github.com/NLnetLabs/unbound

# Update root server hints update every 3 months at 3am on the 1st of the month
0 3 1 */3 * wget -O /etc/unbound/root.hints https://www.internic.net/domain/named.root && service unbound restart